# systemd/user/battery-monitor.service
[Unit]
Description=Battery Monitor

[Service]
ExecStart=%h/.local/bin/battery-monitor.sh
Restart=always

[Install]
WantedBy=default.target
